# settings/ico-dic-wait-weston-ready.path
# Description=Path activation for ico-device-input-controller

[Path]
PathExists=/run/user/5000/wayland-0.lock
#PathExists=/tmp/amd_ready
Unit=ico-device-input-controller.service

[Install]
WantedBy=weston.target 
